# spytti.service
[Unit]
Description=Spytti - Spotify Connect Daemon
After=network-online.target sound.target
Wants=network-online.target sound.target

[Service]
Type=simple
ExecStart=/usr/local/bin/spytti
Restart=always
RestartSec=2
User=root
Group=audio
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spytti

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/cache/spytti
ProtectHome=yes
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
